# source/11 Other Tips and Tricks/11 Step by step a simple data driven application/crontab
#crontab

* *  *  *  * root *applications/bg_processing/cron/process_reports.py
